# gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/favorite-things
command=/home/ubuntu/.local/share/virtualenvs/favorite-things-NiYzl7a6/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/favorite-things/app.sock favorite_things.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[group:guni]
programs:gunicorn